# crontab
# * */6 * * * /usr/local/bin/speedtest.sh
*/5 * * * * /usr/local/bin/speedtest.sh
