# crontab
* */6 * * * /usr/local/bin/speedtest.sh
