# config/crontab
15 */6 * * * /usr/local/bin/speedtest.sh
